# files/user-dir-manager.path
[Unit]
Description=User directory listener

[Install]
WantedBy=multi-user.target

[Path]
PathModified=/home